# irysview_backup.service
[Unit]
Description=Backup IrysView DB on BAS

[Service]
ExecStart=/var/lib/pgsql/pg_irysview.sh
User=postgres
